# .config/autostart/xbindkeys.desktop
[Desktop Entry]
Exec=xbindkeys
Hidden=false;
Icon=
Name=xbindkeys
OnlyShowIn=GNOME;
Path=
StartupNotify=false;
Terminal=False
Type=Application
