# udev/rules.d/60-AC-out.rules
KERNEL=="AC", SUBSYSTEM=="power_supply", ATTR{online}=="0", ATTR{type}=="Mains", RUN+="battery-notify '%k'"
