# udev/rules.d/60-battery-level.rules
KERNEL=="BAT0", SUBSYSTEM=="power_supply", ATTR{capacity}=="[0-19]", RUN+="battery-notify '%k'"
