# udev/rules.d/69-hdparm.rules
ACTION=="add|change", SUBSYSTEM=="block", KERNEL=="sda", RUN+="/usr/bin/hdparm -B 127 /dev/sda"
